# .claude/skills/tokenmax-fleet-orchestrator/assets/systemd/tokenmax-fleet.timer
[Unit]
Description=Run one TokenMax fleet site nightly

[Timer]
OnCalendar=*-*-* 21:00:00 America/Los_Angeles
Persistent=true
RandomizedDelaySec=300
Unit=tokenmax-fleet.service

[Install]
WantedBy=timers.target
